# raspi/image/data.mount.in
[Unit]
Description=DanCam recording data
After=persist.mount
ConditionPathExists=/persist/dancam/storage-admitted

[Mount]
What=/dev/disk/by-label/@DANCAM_DATA_LABEL@
Where=/data
Type=ext4
Options=noatime,errors=remount-ro

[Install]
WantedBy=local-fs.target
