# raspi/dancam.service
# systemd unit for the dancam camera-unit service.
#
# Installed to /etc/systemd/system/dancam.service by raspi/deploy.sh. Auto-starts
# on boot (this is also how the car image auto-records on boot) and restarts on
# crash. Logs go to journald: `journalctl -u dancam -f`.
#
# Footage lives on /data/rec. Do not add RequiresMountsFor=/data here: with the
# playbook's nofail fstab entry, the service must still start so health/status and
# preview remain available when /data is unhealthy. The service-side mount witness
# refuses recording and time-sync mutations instead.

[Unit]
Description=dancam camera-unit service (control + media API)
After=network.target
Wants=network.target

[Service]
Type=exec
# The binary defaults to loopback-only 127.0.0.1:8080. Bind the dual-stack [::]
# wildcard (the binary disables IPV6_V6ONLY) so the phone on the AP (IPv4) and
# the Mac on the LAN (IPv4 or IPv6) can reach it. 8080 is unprivileged, so no
# root / CAP_NET_BIND_SERVICE needed.
Environment=DANCAM_BIND=[::]:8080
# The deployed Pi uses the real camera backend. Local cargo runs default to the
# mock backend unless DANCAM_BACKEND=camera is set explicitly.
Environment=DANCAM_BACKEND=camera
# Recording directory on the dedicated data partition. The mount witness prevents
# silently creating /data/rec on the root filesystem if /data is absent.
Environment=DANCAM_REC_DIR=/data/rec
Environment=DANCAM_REQUIRE_REC_MOUNT=/data
ExecStart=/usr/local/bin/dancam
Restart=on-failure
RestartSec=2
User=dancam
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
